# root/etc/s6-overlay/s6-rc.d/p21-rpigpio-customize/up
/etc/s6-overlay/s6-rc.d/p21-rpigpio-customize/run
