# root/etc/s6-overlay/s6-rc.d/p20-rpigpio-setup/up
/etc/s6-overlay/s6-rc.d/p20-rpigpio-setup/run
